# supervisor/chat.conf
[program:chat]
directory=/home/hello/chat/
command=/home/hello/chat/chat.sh
user=hello
autorestart=true
